# Course5/Lab3/hello_cloud.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/hello_cloud.py 80

[Install]
WantedBy=default.target
